# override/etc/services.d/20-nginx/run
#!/usr/bin/with-contenv sh
exec nginx -c /mnt/config/etc/nginx.conf;