# override/etc/s6-overlay/s6-rc.d/certbot/up
#!/command/execlineb
# SPDX-FileCopyrightText: 2022 - 2024 Daniel Wolf <nephatrine@gmail.com>
#
# SPDX-License-Identifier: ISC

foreground { /command/with-contenv /bin/bash /usr/local/bin/certify }
